# uliweb/template_files/config/supervisor.conf
[program:{{=project}}]
command = uwsgi
 --socket /tmp/{{=project}}.sock
 --harakiri 60
 --reaper
 --module wsgi_handler
 --processes 10
 --master{{if defined('python_path'):}}{{out.write(" --home %s\n" % python_path)}}{{pass}}
 --chmod-socket=666
 --limit-as 512
 --socket-timeout 5
 --max-requests 1000
 --reload-mercy 8
 --reload-on-as 384
 --reload-on-rss 192
 --no-orphans
directory={{=project_dir}}
stopsignal=QUIT
autostart=true
autorestart=true
stdout_logfile=/tmp/{{=project}}.log
redirect_stderr=true
exitcodes=0,1,2
